# openrc/servicio-db-blog
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
    need net
    need localmount
    need nfsmount
    after sshd
}

start() {
	ebegin "Arrancando servicios: BD MySQL para Blog Wordpress"

        su luis -c "cd /Apps/docker/servicio-db-blog && /usr/local/bin/fig up -d"

        eend $?
}

stop() {
	ebegin "Parando servicios: BD MySQL para Blog Wordpress"

        su luis -c "cd /Apps/docker/servicio-db-blog &&  /usr/local/bin/fig stop"

	eend $?
}

status() {

        su luis -c "cd /Apps/docker/servicio-db-blog &&  /usr/local/bin/fig ps"

	eend $?
}
